# vendor/l3-transport/launchagents/com.praise.codex-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Placeholders (__HOME__ / __HANDOFF_BIN__ / __HANDOFF_STATE__ / __CLAUDE_PLUGIN_DATA__)
     are rendered by install-transport.sh against state-paths.sh-resolved values at deploy time.
     This bundle file MUST NOT contain any developer-machine absolute path.
     __CLAUDE_PLUGIN_DATA__ is required so the daemon resolves the SAME gd-handoff root as
     the client (state-paths.sh fallback); without it daemon/client drift to different paths. -->
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>__HOME__</string>
		<key>PATH</key>
		<string>__HOME__/.npm-global/bin:__HOME__/.bun/bin:__HOME__/.pyenv/shims:__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
		<key>CODEX_PROVIDER_CHAIN</key>
		<string>tapsvc</string>
		<key>CODEX_RETRY_BACKOFF</key>
		<string>0</string>
		<key>CODEX_SEND_WAIT_TIMEOUT</key>
		<string>1500</string>
		<key>CODEX_EXEC_TIMEOUT</key>
		<string>720</string>
		<key>CLAUDE_PLUGIN_DATA</key>
		<string>__CLAUDE_PLUGIN_DATA__</string>
		<key>TAPTAP_API_KEY</key>
		<string>__TAPTAP_API_KEY__</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.praise.codex-watch</string>
	<key>ProgramArguments</key>
	<array>
		<string>__HANDOFF_BIN__/codex-watch</string>
		<string>run</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>__HANDOFF_STATE__/codex-watch-launchd.log</string>
	<key>StandardOutPath</key>
	<string>__HANDOFF_STATE__/codex-watch-launchd.log</string>
	<key>ThrottleInterval</key>
	<integer>5</integer>
</dict>
</plist>
