# vendor/l3-transport/launchagents/com.praise.codex-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>/Users/praise</string>
		<key>PATH</key>
		<string>/Users/praise/.npm-global/bin:/Users/praise/.bun/bin:/Users/praise/.pyenv/shims:/Users/praise/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
		<key>TAPTAP_API_KEY</key>
		<string>&lt;REDACTED-set-via-install-or-launchctl-setenv&gt;</string>
		<key>CODEX_PROVIDER_CHAIN</key>
		<string>tapsvc</string>
		<key>CODEX_RETRY_BACKOFF</key>
		<string>0</string>
		<key>CODEX_SEND_WAIT_TIMEOUT</key>
		<string>540</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.praise.codex-watch</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/praise/.claude/handoff/bin/codex-watch</string>
		<string>run</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/praise/.claude/handoff/state/codex-watch-launchd.log</string>
	<key>StandardOutPath</key>
	<string>/Users/praise/.claude/handoff/state/codex-watch-launchd.log</string>
	<key>ThrottleInterval</key>
	<integer>5</integer>
</dict>
</plist>